# dotfiles/dot_config/systemd/user/battery-monitor.timer
[Unit]
Description=Check battery level every 30 seconds

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
